# projects/magmadar/files/supervisord.d/magmadar-ng.ini
[program:magmadar-ng]
process_name=%(program_name)s
command=/search/srv/magmadar-ng/magmadar/magmadar -p 18361 -c online
numprocs=1
directory=/search/srv/magmadar-ng/magmadar
autostart=false
startretries=1
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/s/logs/magmadar/magmadar-ng.out
stderr_logfile=/s/logs/magmadar/magmadar-ng.err
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=6


[program:magmadar-ng-backup]
process_name=%(program_name)s
command=/search/srv/magmadar-ng/magmadar-backup/magmadar -p 18362 -c online
numprocs=1
directory=/search/srv/magmadar-ng/magmadar-backup
autostart=false
startretries=1
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/s/logs/magmadar/magmadar-ng.out
stderr_logfile=/s/logs/magmadar/magmadar-ng.err
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=6
