# projects/magmadar/files/supervisord.d/magmadar.ini
[program:magmadar]
process_name=%(program_name)s-%(process_num)d
command=/search/srv/magmadar/venv/magmadar/bin/python -m magmadar.web --socket=/tmp/magmadar_sockets/sock.%(process_num)d
environment=MAGMADAR_CONF=online.cfg
numprocs=20
directory=/search/srv/magmadar/magmadar
autostart=false
startretries=1
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/s/logs/magmadar/magmadar.out.%(process_num)d
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=6
stderr_logfile=/s/logs/magmadar/magmadar.err.%(process_num)d
stderr_logfile_maxbytes=500MB
stderr_logfile_backups=6
