# projects/magmadar/files/supervisord.d/redis.ini
[program:magmadar-redis]
process_name=%(program_name)s-%(process_num)d
command=/search/local/bin/redis-server /search/srv/redis/magmadar.conf --port %(process_num)d
numprocs=8
numprocs_start=16380
directory=/search/srv/redis
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/s/logs/magmadar/redis.out.%(process_num)d
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=3
stderr_logfile=/s/logs/magmadar/redis.err.%(process_num)d
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=3
