# projects/magmadar/files/supervisord.d/twemproxy.ini
[program:magmadar-twemproxy]
process_name=%(program_name)s-%(process_num)d
command=/search/srv/twemproxy/nutcracker -c /search/srv/twemproxy/magmadar.yml.%(process_num)d -v 4 -m 512 -s %(process_num)d
numprocs=2
numprocs_start=22222
directory=/search/srv/twemproxy
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/s/logs/magmadar/twemproxy.out
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=3
stderr_logfile=/s/logs/magmadar/twemproxy.err
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=3
